# h_to_rocket_chat/h_to_rocket_chat.supervisor.conf.example
[program:h-to-rocket-chat]
command = <path_to_h_to_rocket_chat_py> --config <path_to_rocket_chat_conf>
